# web-app/miner-v2/imac/com.razum.miner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.razum.miner</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/svatoslav/.razum/miner.js</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin</string>
    <key>MODELS_CSV</key>
    <string>qwen2.5:7b</string>
    <key>MINER_NAME</key>
    <string>imac-svatoslav</string>
    <key>WALLET_ADDRESS</key>
    <string>0x60769DFDF421d6633c0b0a571936272c93A63Db2</string>
  </dict>
  <key>WorkingDirectory</key>
  <string>/Users/svatoslav/.razum</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/svatoslav/.razum/miner.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/svatoslav/.razum/miner.err.log</string>
</dict>
</plist>
